# roles/powertop/templates/powertop.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Extend the battery life of laptop

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/powertop --auto-tune
{% if ansible_facts.hostname == 'bagel' %}
# Following line is a workaround for ethernet card interruptions issue
ExecStart=/bin/sh -c "echo 'on' > /sys/bus/pci/devices/0000:00:1c.6/power/control"
{% endif %}

[Install]
WantedBy=multi-user.target
